# host/systemd/clawfarmer-host-photo-evening.timer
[Unit]
Description=clawfarmer evening photo (19:00 local)

[Timer]
OnCalendar=*-*-* 19:00:00
Unit=clawfarmer-host-tick@photo.service
Persistent=true

[Install]
WantedBy=timers.target
