# host/systemd/clawfarmer-host-sensors.timer
[Unit]
Description=clawfarmer sensor sweep every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Unit=clawfarmer-host-tick@sensors.service
Persistent=true
RandomizedDelaySec=30s

[Install]
WantedBy=timers.target
