# host/systemd/clawfarmer-host-tick@.service
[Unit]
Description=clawfarmer host tick (%I)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=openclaw
Group=openclaw
Environment=HOME=/var/lib/openclaw
Environment=CLAWFARMER_WORKSPACE=/var/lib/openclaw/.openclaw/workspace-plant
Environment=PI_HOST=clawpi.local
Environment=PI_USER=pi
Environment=PI_KEY=/var/lib/openclaw/.ssh/id_ed25519_plantpi
Environment=JETSON_HOST=orin-nano.local
Environment=JETSON_USER={{JETSON_USER}}
Environment=JETSON_KEY=/var/lib/openclaw/.ssh/id_ed25519_plantjetson
Environment=SOIL_CHANNEL=0
Environment=SOIL_DRY_RAW=26000
Environment=SOIL_WET_RAW=12000
# Optional — if set, each successful photo capture triggers this cron job
# (used to refresh the rich photo-review analysis on every capture).
Environment=PHOTO_REVIEW_CRON_ID={{PHOTO_REVIEW_CRON_ID}}
# Set to "1" to run rich analysis on the Jetson (Gemma) instead of via the
# OpenClaw cron chain (Kimi/Qwen). Saves Together API cost but Telegram
# stops getting rich plant-check messages (dashboard still updates).
Environment=USE_JETSON_RICH_ANALYSIS=0
ExecStart=/usr/local/bin/clawfarmer-host-tick %I
StandardOutput=journal
StandardError=journal
