# host/systemd/clawfarmer-host-photo-every-4h.timer
[Unit]
Description=clawfarmer photo every 4 hours (daytime only, 08–20 local)

[Timer]
OnCalendar=*-*-* 08..20/4:00:00
Unit=clawfarmer-host-tick@photo.service
Persistent=true
RandomizedDelaySec=30s

[Install]
WantedBy=timers.target
