# host/systemd/clawfarmer-openclaw-restart.timer
[Unit]
Description=Restart openclaw daily at 03:00 to prevent memory leak accumulation

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
RandomizedDelaySec=2m

[Install]
WantedBy=timers.target
